# fangfrisch/config/crontab
# Run fangfrisch every 10 minutes
*/10 * * * * /app/venv/bin/fangfrisch --conf /etc/fangfrisch.conf refresh
